# devops/03/slice.slice
[Unit]
Description=Slice


[Slice]
MemoryMax=500M
CPUQuota=20%
