# devops/03/z6/slice.slice
[Unit]
Description=Slice
Documentation=man:systemd.slice(5)

[Slice]
MemoryMax=200M 
CPUQuota=20%
